# config.qemu/systemd/sway-wallpaper.service
[Unit]
Description=Sway Wallpaper Daemon (runs 24/7)
# Note: This service is now started automatically by Sway config
# It's kept here for manual control if needed

[Service]
Type=simple
ExecStart=%h/.local/bin/sway-wallpaper-daemon.sh
Restart=always
RestartSec=10
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
# Disabled auto-start - daemon is launched by Sway config instead
# WantedBy=default.target
