# config.mac/systemd/package-sync.service
[Unit]
Description=Package Sync Daemon (runs 24/7)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/package-sync-daemon.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target