# config.mac/systemd/sway-theme-switcher.service
[Unit]
Description=Sway Theme Switcher Daemon (runs 24/7)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/sway-theme-switcher-daemon.sh
Restart=always
RestartSec=10
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
